# deploy/systemd/biascheck-backend.service
[Unit]
Description=Bias Checker FastAPI backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/biasCheck/backend
EnvironmentFile=/var/www/biasCheck/backend/.env
ExecStart=/var/www/biasCheck/backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
